# templates/webalizer.cron.erb
#!/bin/bash

# File installed with puppet.

if [ -r <%= @logfile %> ] ; then
  /usr/bin/webalizer -Q
fi

for CONFIG in /etc/webalizer/*.conf; do
  /usr/bin/webalizer -c ${CONFIG}
done
